# examples/configs/com.user.dashboard-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.dashboard-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/<user>/bin/dashboard-sync.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>REMOTE</key>
        <string>linux</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/<user>/.dashboard-sync.stdout</string>
    <key>StandardErrorPath</key>
    <string>/Users/<user>/.dashboard-sync.stderr</string>
</dict>
</plist>
