# examples/configs/mini-router.service.example
# systemd user unit — mini-router (Telegram → tmux forwarder)
# 위치: ~/.config/systemd/user/mini-router.service
# 활성화: systemctl --user daemon-reload && systemctl --user enable --now mini-router.service
# 로그: journalctl --user -u mini-router -f  (또는 ~/.hermes/logs/mini-router.stderr.log)
#
# 사전 준비:
#   1. /home/YOUR_USERNAME/agentic-harness/mini-router 안에 .venv 생성 후
#      pip install -r requirements.txt (python-telegram-bot 21.x + PyYAML)
#   2. secret-tool로 telegram-bot-token 저장 (attribute 순서 주의 — service hermes account telegram-bot-token):
#        printf '%s' '<BOT_TOKEN>' | secret-tool store --label='Telegram bot' service hermes account telegram-bot-token
#   3. ~/.hermes/.env에 user_id / chat_id 작성 (chmod 600)
#      — EnvironmentFile=- 라서 파일이 없어도 기동 자체는 막지 않음 (토큰은 secret-tool에서 lookup)
#      — 주의: systemd EnvironmentFile은 $HOME 변수를 확장하지 않으므로 절대 경로 권장
#   4. mkdir -p ~/.hermes/logs (StandardOutput append: 대상 디렉토리가 미리 있어야 함)
#   5. loginctl enable-linger $USER (부팅 시 자동 시작)
#
# 경로 주의: 아래 /home/YOUR_USERNAME/agentic-harness 는 placeholder입니다.
#            install 스크립트가 실제 클론 경로로 치환하거나, 수동 설치 시 직접 수정하세요.

[Unit]
Description=Telegram → tmux/opencode mini-router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USERNAME/agentic-harness/mini-router
EnvironmentFile=-/home/YOUR_USERNAME/.hermes/.env

# secret-tool에서 토큰 lookup해서 환경변수로 export 후 python 실행
ExecStart=/bin/bash -c 'export TELEGRAM_BOT_TOKEN=$(secret-tool lookup service hermes account telegram-bot-token) && exec /home/YOUR_USERNAME/agentic-harness/mini-router/.venv/bin/python /home/YOUR_USERNAME/agentic-harness/mini-router/bot.py'

Restart=on-failure
RestartSec=10

StandardOutput=append:/home/YOUR_USERNAME/.hermes/logs/mini-router.stdout.log
StandardError=append:/home/YOUR_USERNAME/.hermes/logs/mini-router.stderr.log

# PATH에 opencode 포함 (tmux send-keys 시 자동 부팅에 필요)
Environment=PATH=/home/YOUR_USERNAME/.opencode/bin:/home/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin

# 라우팅 기본값
Environment=WORKDIR_ROOT=/home/YOUR_USERNAME/dev
Environment=TMUX_DEFAULT_SESSION=oc-default

[Install]
WantedBy=default.target
